# server/update_ods.cron
0 * * * * bash /home/ubuntu/docker_toto/server/update_redis.py

